# apt/iot4bd/DEBIAN/postrm
#!/bin/bash
set -e

case "$1" in
    remove|purge)
        systemctl disable iot4bd.service || true
        systemctl daemon-reload || true
        ;;
esac
