# iot4b_package/etc/init.d/iot4b
#!/bin/sh /etc/rc.common
START=50
STOP=50

start() {
    /opt/iot4b/iot4b -env=iot4b &
}

stop() {
    killall iot4b
}