# apt/iot4b/DEBIAN/postrm
#!/bin/bash
set -e

case "$1" in
    remove|purge)
        systemctl disable iot4b.service || true
        systemctl daemon-reload || true
        ;;
esac
